# scepg-wget/scepg-wget.service
[Unit]
Description=scepg-wget

[Service]
Type=oneshot
WorkingDirectory=/var/lib/scepg
ExecStart=/var/lib/scepg/scepg-wget
