# epg2xml-wget/epg2xml-wget.service
[Unit]
Description=epg2xml-wget

[Service]
Type=oneshot
WorkingDirectory=/var/lib/epg2xml
ExecStart=/var/lib/epg2xml/epg2xml-wget
